# kill_papers.cron
30	18	*	*	*	root	/usr/bin/python /usr/sbin/kill_papers.py
